# systemd/agent-control-surface.service
[Unit]
Description=agent-control-surface (local-only)
After=network-online.target

[Service]
Type=simple
WorkingDirectory=%h/repos/heimgewebe/agent-control-surface
ExecStart=%h/repos/heimgewebe/agent-control-surface/.venv/bin/uvicorn panel.app:app --host 127.0.0.1 --port 8099
Restart=always
RestartSec=2

[Install]
WantedBy=default.target
